# cron.d
01 04 * * * /home/userName/backuputils/mysql-backupall.sh ~/dbdump
02 04 * * * /home/userName/backuputils/mysql-cleanupall.sh 5 /home/userName/dirWithDumps
