# cron.daily
01 04 * * * /home/$USER/backuputils/mysql-backupall.sh ~/dbdump
02 04 * * * /home/$USER/backuputils/mysql-cleanupall.sh 5 /home/$USER/dirWithDumps
